# var/lib/dpkg/info/alsa-utils.list
/.
/usr
/usr/sbin
/usr/sbin/alsactl
/usr/sbin/alsaconf
/usr/share
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/amidi.1.gz
/usr/share/man/man1/amixer.1.gz
/usr/share/man/man1/aplay.1.gz
/usr/share/man/man1/iecset.1.gz
/usr/share/man/man1/aplaymidi.1.gz
/usr/share/man/man1/arecordmidi.1.gz
/usr/share/man/man1/aseqdump.1.gz
/usr/share/man/man1/aseqnet.1.gz
/usr/share/man/man1/speaker-test.1.gz
/usr/share/man/man1/asoundconf.1.gz
/usr/share/man/man1/alsactl.1.gz
/usr/share/man/man1/alsamixer.1.gz
/usr/share/man/man1/aconnect.1.gz
/usr/share/man/man8
/usr/share/man/man8/alsaconf.8.gz
/usr/share/locale
/usr/share/locale/ja
/usr/share/locale/ja/LC_MESSAGES
/usr/share/locale/ja/LC_MESSAGES/alsaconf.mo
/usr/share/locale/ja/LC_MESSAGES/alsa-utils.mo
/usr/share/locale/ru
/usr/share/locale/ru/LC_MESSAGES
/usr/share/locale/ru/LC_MESSAGES/alsaconf.mo
/usr/share/alsa
/usr/share/alsa/speaker-test
/usr/share/alsa/speaker-test/sample_map.csv
/usr/share/sounds
/usr/share/sounds/alsa
/usr/share/sounds/alsa/Front_Left.wav
/usr/share/sounds/alsa/Rear_Center.wav
/usr/share/sounds/alsa/Rear_Right.wav
/usr/share/sounds/alsa/Side_Right.wav
/usr/share/sounds/alsa/Front_Center.wav
/usr/share/sounds/alsa/Front_Right.wav
/usr/share/sounds/alsa/Noise.wav
/usr/share/sounds/alsa/Rear_Left.wav
/usr/share/sounds/alsa/Side_Left.wav
/usr/share/doc
/usr/share/doc/alsa-utils
/usr/share/doc/alsa-utils/README.Debian
/usr/share/doc/alsa-utils/copyright
/usr/share/doc/alsa-utils/PATCHES.Debian
/usr/share/doc/alsa-utils/README.alsamixer
/usr/share/doc/alsa-utils/changelog.gz
/usr/share/doc/alsa-utils/changelog.Debian.gz
/usr/share/menu
/usr/share/menu/alsa-utils
/usr/bin
/usr/bin/alsamixer
/usr/bin/amidi
/usr/bin/amixer
/usr/bin/aplay
/usr/bin/iecset
/usr/bin/aconnect
/usr/bin/aplaymidi
/usr/bin/arecordmidi
/usr/bin/aseqdump
/usr/bin/aseqnet
/usr/bin/speaker-test
/usr/bin/asoundconf
/etc
/etc/alsa
/etc/alsa/modprobe-post-install.d
/etc/alsa/modprobe-post-install.d/alsa-utils
/etc/udev
/etc/udev/rules.d
/etc/udev/alsa-utils.rules
/etc/init.d
/etc/init.d/alsa-utils
/lib
/lib/udev
/lib/udev/alsa-utils
/var
/var/lib
/var/lib/alsa
/usr/share/man/man1/arecord.1.gz
/usr/bin/arecord
