# var/lib/dpkg/info/ifplugd.conffiles
/etc/ifplugd/action.d/ifupdown
/etc/ifplugd/ifplugd.action
/etc/init.d/ifplugd
/etc/apm/scripts.d/ifplugd
/etc/udev/rules.d/z60_ifplugd.rules
