# var/lib/dpkg/info/php5-common.list
/etc
/etc/cron.d
/etc/cron.d/php5
/etc/php5
/etc/php5/conf.d
/etc/php5/conf.d/pdo.ini
