# var/lib/dpkg/info/udev.conffiles
/etc/udev/rules.d/95-late.rules
/etc/udev/rules.d/60-persistent-storage-tape.rules
/etc/udev/rules.d/91-permissions.rules
/etc/udev/rules.d/75-cd-aliases-generator.rules
/etc/udev/rules.d/60-persistent-v4l.rules
/etc/udev/rules.d/50-udev.rules
/etc/udev/rules.d/75-persistent-net-generator.rules
/etc/udev/rules.d/60-persistent-storage.rules
/etc/udev/rules.d/60-persistent-input.rules
/etc/udev/rules.d/80-drivers.rules
/etc/udev/udev.conf
/etc/udev/links.conf
/etc/scsi_id.config
/etc/init.d/udev
/etc/init.d/udev-mtab
/etc/modprobe.d/display_class
/etc/modprobe.d/blacklist
/etc/modprobe.d/pnp-hotplug
