# var/lib/dpkg/info/udhcpd.list
/.
/var
/var/lib
/var/lib/misc
/var/lib/misc/udhcpd.leases
/etc
/etc/init.d
/etc/init.d/udhcpd
/usr
/usr/bin
/usr/bin/dumpleases
/usr/share
/usr/share/doc
/usr/share/doc/udhcpd
/usr/share/doc/udhcpd/changelog.Debian.gz
/usr/share/doc/udhcpd/changelog.gz
/usr/share/doc/udhcpd/README.gz
/usr/share/doc/udhcpd/examples
/usr/share/doc/udhcpd/examples/udhcpd.conf
/usr/share/doc/udhcpd/copyright
/usr/share/man
/usr/share/man/man5
/usr/share/man/man5/udhcpd.conf.5.gz
/usr/share/man/man1
/usr/share/man/man1/dumpleases.1.gz
/usr/share/man/man8
/usr/share/man/man8/udhcpd.8.gz
/usr/sbin
/usr/sbin/udhcpd
